# app/tags/1.0.1/etc/nddi-scheduler.cron
* * * * * root /usr/bin/perl /usr/bin/oess_scheduler.pl
1 * * * * root /usr/bin/perl /usr/bin/snapp-config-gen