# app/trunk/etc/nddi-scheduler.cron
* * * * * root /usr/bin/perl /usr/bin/oess_scheduler.pl
1 1 * * * root /usr/bin/perl /usr/bin/populate_remote_information.pl
